# setup/relay-etiquetas.service
[Unit]
Description=Relay Etiquetas de Caducidad (Flask + gunicorn)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/home/n1ce/relay-etiquetas
# Puerto 80 requiere root (o usar setcap / authbind). gunicorn con 2 workers.
ExecStart=/home/n1ce/relay-etiquetas/venv/bin/gunicorn \
    --workers 2 \
    --bind 0.0.0.0:80 \
    app:app
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
